# lxqt/dot-config/autostart/Autostart Cleanup.desktop
[Desktop Entry]
Exec=/home/deorz/.config/hypr/scripts/cleanup.sh
Name=Autostart Cleanup
OnlyShowIn=LXQt;
Type=Application
Version=1.0
