# lxqt/dot-config/autostart/Clipboard Manager.desktop
[Desktop Entry]
Exec=wl-paste --watch cliphist store
Name=Clipboard Manager
OnlyShowIn=LXQt;
Type=Application
Version=1.0
