# lxqt/dot-config/autostart/NWG Panel.desktop
[Desktop Entry]
Exec=nwg-panel
Hidden=true
Name=NWG Panel
OnlyShowIn=LXQt;
Type=Application
Version=1.0
