# lxqt/dot-config/autostart/Wallpaper Daemon.desktop
[Desktop Entry]
Exec=swww-daemon
Name=Wallpaper Daemon
OnlyShowIn=LXQt;
Type=Application
Version=1.0
