# lxqt/dot-config/autostart/Waybar.desktop
[Desktop Entry]
Exec=waybar
Name=Waybar
OnlyShowIn=LXQt;
Type=Application
Version=1.0
